# systemd/nlweb-sidecar.service
[Unit]
Description=NLWeb Sidecar (aiohttp)
After=network.target

[Service]
Type=simple
User=%i
EnvironmentFile=/etc/esw_know.env
WorkingDirectory=/opt/esw.know
ExecStart=/usr/bin/env python3 /opt/esw.know/NLWeb/webserver/aiohttp_server.py
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
[Unit]
Description=NLWeb Sidecar (aiohttp)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/workspaces/esw.know
EnvironmentFile=/etc/esw_know.env
ExecStart=/usr/bin/env python3 /workspaces/esw.know/NLWeb/webserver/aiohttp_server.py --host 127.0.0.1 --port 7015
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
